# 082619_1/one-daemon.service
[Unit]
Description="1st daemon"

[Service]
ExecStart=/tmp/one-daemon

