# 082619_1/two-daemon.timer
[Unit]
Description="2nd daemon timer"

[Timer]
OnCalendar=2019-01-01
