# 082619_1/two-daemon.service
[Unit]
Description="2nd daemon"
Requires=one-daemon.service
After=one-daemon.service

[Service]
Type=oneshot
ExecStart=/tmp/two-daemon

