# 092319/exam22.service
[Unit]
Description=Exam22
After=network.target

[Service]
Type=oneshot
ExecStart=/tmp/exam/exam22.cmd

[Install]
WantedBy=multi-user.target
