# 090519/ext-daemon.service
[Unit]
Description="User Verification"

[Service]
ExecStart=/tmp/ext-daemon.cmd
